# swift-files/common/storage/drive-audit.cron
*/5 * * * * root /usr/bin/swift-drive-audit /etc/swift/drive-audit.conf 2>&1 | logger
